# drunkstar/skeleton/etc/network/if-up.d/ethtool
#!/bin/sh

ethtool -s eth0 autoneg off duplex full speed 10
